# scripts/miniros-high.slice
[Unit]
Description=Slice for MiniROS services with high performance reqirements

[Slice]

[Install]
WantedBy=multi-user.target